# deploy/supervisor/kis-scheduler.conf
; KISA Laravel scheduler.
;
; Runs the Laravel scheduler as a long-lived process (schedule:work ticks once
; per minute), which is what actually fires the scheduled commands defined in
; routes/console.php — stats:generate, the Billbee syncs, and the nightly
; backup. Without this, those commands never run on their own.
;
; Install: copy to /etc/supervisor/conf.d/kis-scheduler.conf on the server, then
;   sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start kis-scheduler
;
; Runs as www-data to match kis-octane, so it shares the app user's Redis,
; cache, and log access.
[program:kis-scheduler]
directory=/var/www/kis.kraeuter-wege.de
command=php artisan schedule:work
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/kis.kraeuter-wege.de/storage/logs/scheduler.log
